# .svn/pristine/36/362e891bb56d5cb644ee414e8f93007b7ec41add.svn-base
#!/bin/bash
# description: MPS Platform Start, Primary Environment
# processname: mps
# chkconfig: - 89 11
CAR_HOME=/home/mpsadmin/car

case $1 in
start)
/bin/su mpsadmin -c 'python /home/mpsadmin/car/commands/skynyrd.py --dst /usr/local/mps'
/bin/su mpsadmin -c 'python /home/mpsadmin/car/MPSAuthSvc.py -e devserver --file config/devserver/MPSAuthSvc/config.json &'
/bin/su mpsadmin -c 'python /home/mpsadmin/car/MPSAdmin.py -e devserver --file config/devserver/MPSAdmin/config.json &'
/bin/su mpsadmin -c 'python /home/mpsadmin/car/MPSCV.py -e devserver --file config/devserver/MPSCV/config.json &'
/bin/su mpsadmin -c 'python /home/mpsadmin/car/MPSCV.py -e devserver --file config/devserver/MPSCV/config_a.json &'
/bin/su mpsadmin -c 'python /home/mpsadmin/car/MPSAppt.py -e devserver --file config/devserver/MPSAppt/config.json &'
/bin/su mpsadmin -c 'python /home/mpsadmin/car/MPSAppt.py -e devserver --file config/devserver/MPSAppt/config_a.json &'
/bin/su mpsadmin -c 'python /home/mpsadmin/car/MPSLogin.py -e devserver --file config/devserver/MPSLogin/config.json &'
;;
esac
exit 0
